# old/airootfs/usr/lib/systemd/system/conky.service
[Unit]
Description=conky fix

[Service]
User=root
Type=oneshot
ExecStart=/etc/systemd/scripts/conky

[Install]
WantedBy=multi-user.target

