# old/airootfs/root/systemd/pinentry.service
[Unit]
Description=pinentry fix

[Service]
User=root
Type=oneshot
ExecStart=/etc/systemd/scripts/pinentry

[Install]
WantedBy=multi-user.target

