# old/grsec.archiso/airootfs/usr/lib/systemd/system/dhcpd4.service
[Unit]
Description=IPv4 DHCP server
After=network.target

[Service]
#User=dhcp
Type=forking
PIDFile=/run/dhcpd4.pid
ExecStart=/usr/bin/dhcpd -4 -q -pf /run/dhcpd4.pid internal
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
