# old/grsec.archiso/airootfs/usr/lib/systemd/system/macchanger.service
[Unit]
Description=Mac Spoof

[Service]
Type=oneshot
ExecStart=/etc/systemd/scripts/macchanger

[Install]
WantedBy=network.target

