# old/grsec.archiso/airootfs/usr/lib/systemd/system/nat-bridge.service
[Unit]
Description=nat bridge

[Service]
Type=oneshot
ExecStart=/etc/systemd/scripts/nat-bridge

[Install]
WantedBy=network.target

