# old/minimal.archiso/airootfs/etc/systemd/system/getty@tty3.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin admin --noclear %I 38400 linux
