# airootfs/usr/lib/systemd/system/fcgiwrap.socket
[Unit]
Description=fcgiwrap Socket

[Socket]
ListenStream=/home/nginx/fcgiwrap.sock

[Install]
WantedBy=sockets.target
