# airootfs/usr/lib/systemd/system/sshd-spareproject.service
[Unit]
Description=OpenSSH Daemon
After=network.target

[Service]
User=spareproject
ExecStart=/usr/bin/sshd -D -f /home/spareproject/sshd/sshd_config -E /home/spareproject/sshd/log
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target

# This service file runs an SSH daemon that forks for each incoming connection.
# If you prefer to spawn on-demand daemons, use sshd.socket and sshd@.service.
